# rhel-jetson-triton/files/etc/systemd/system/mac-hostname.service
[Unit]
Description=Set hostname to MAC address
Documentation=man:hostnamectl(1)
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set-hostname-from-mac.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
Environment=HOSTNAME_INTERFACE=enP8p1s0

# Restart policy
Restart=no

# Security settings
PrivateNetwork=false
PrivateDevices=false
ProtectSystem=false
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target